# .config/systemd/user/mopidy.service
[Unit]
Description=Mopidy music server
After=network-online.target
Wants=network-online.target
After=pulseaudio.service
After=sound.target

[Service]
ExecStart=/usr/bin/mopidy --config /home/jtao/.config/mopidy/mopidy.conf

[Install]
WantedBy=default.target
